# sshd.conf
[program:ss-wrapper]
priority = 9
command = /usr/sbin/sshd -D
autostart = true
startsecs = 0
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/sshd
stdout_events_enabled = true